# docker/root/etc/s6-overlay/s6-rc.d/tixati/run
#!/command/execlineb -P

with-contenv

backtick -E HOME { homeof user }
s6-env HOME="$HOME"

s6-setuidgid user tixati %U
