# docker/root/etc/s6-overlay/s6-rc.d/tixati/run
#!/command/execlineb -P

with-contenv

tixati %U